# notifications.ini
[program:websocket_server]
autorestart=true
autostart=true
command=python3 /code/websocket_server.py --port %(ENV_NOTIFICATIONS_PORT)s --host %(ENV_NOTIFICATIONS_HOST)s
killasgroup=true
numprocs=1
priority=1
redirect_stderr=true
stopwaitsecs=5
stopasgroup=true

[program:consumer]
autorestart=true
autostart=true
command=python3 /code/consumer.py
numprocs=1
priority=5
redirect_stderr=true

[group:notifications]
programs=consumer,websocket_server
